# update-mirrors.service
[Service]
ExecStart=/bin/bash /usr/local/bin/update-mirrors.sh
Type=oneshot

[Unit]
Requires=network.target
After=network.target
